# debian/splitdnsd-selfbuild/lib/systemd/system/splitdnsd-autobuild.service
[Unit]
Description=Rebuild splitdnsd against the current apt Go toolchain (unattended security self-patch)
Documentation=man:splitdnsd(8)
After=network-online.target
Wants=network-online.target
# Off the apt transaction, so a build failure can never break apt or DNS. On failure,
# email the administrator (the work is asynchronous, so apt does not report it).
OnFailure=splitdnsd-autobuild-notify@%n.service

[Service]
Type=oneshot
ExecStart=/usr/lib/splitdnsd-selfbuild/splitdnsd-autobuild
# Be a good citizen on a small resolver: low priority, idle I/O, modest CPU share. The
# rebuild only compiles (it does not run the test suite — the shipped source is already
# CI-tested), so memory stays modest; GOMAXPROCS is capped in the script/config.
Nice=15
IOSchedulingClass=idle
CPUWeight=20
TimeoutStartSec=30min
# Reuse the cache dir across runs (build artifacts, rollback snapshots).
CacheDirectory=splitdnsd-selfbuild
