# packaging/autobuild/splitdnsd-autobuild.timer
[Unit]
Description=Weekly backstop: rebuild splitdnsd against the apt Go toolchain if it changed
Documentation=man:splitdnsd(8)

[Timer]
# The apt hook (80splitdnsd-autobuild) is the primary trigger — it fires right after the
# toolchain is upgraded. This weekly run is a backstop in case an apt upgrade happened in
# a way that bypassed the hook, or the resolver was offline during one. It is a cheap
# no-op unless the toolchain actually moved.
OnCalendar=weekly
RandomizedDelaySec=6h
Persistent=true

[Install]
WantedBy=timers.target
